# internal/golden/guest/io.goodkind.gha-mac-broker-guest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.goodkind.gha-mac-broker-guest</string>
  <!-- Run as the runner user (admin), not root. The provisioner installs the
       actions runner under this user's home (build.go resolves it via printenv
       HOME, which runs as admin), and the runner refuses to run as root, so the
       guest agent must run as admin with HOME set. A root LaunchDaemon leaves
       HOME empty and every RunJob fails with "base home is empty". -->
  <key>UserName</key>
  <string>admin</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/gha-mac-broker</string>
    <string>guest-supervisor</string>
  </array>
  <!-- The runner inherits this process's environment (guestexec mergedEnvironment
       merges os.Environ() with per-slot overrides), and a LaunchDaemon starts with
       only the minimal launchd PATH (/usr/bin:/bin:/usr/sbin:/sbin). That finds
       system swift/xcodebuild/git but NOT /opt/homebrew/bin, so a job that runs
       brew (e.g. installing swift-format/swiftlint) fails with exit 127. Set the
       runner user's login-shell PATH so brew-installed and rbenv-shimmed tools
       resolve the same way they do in a developer login session. -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/admin</string>
    <key>PATH</key>
    <string>/Users/admin/.rbenv/shims:/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/System/Cryptexes/App/usr/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>KeepAlive</key>
  <true/>
  <key>RunAtLoad</key>
  <true/>
  <!-- /var/log is root-only; admin writes the log to a world-writable path. -->
  <key>StandardOutPath</key>
  <string>/tmp/gha-mac-broker-guest.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/gha-mac-broker-guest.log</string>
</dict>
</plist>
